# deploy/dashboard-workers.conf
[program:dashboard-workers]
process_name=%(program_name)s_%(process_num)02d
#command=php /home/forge/app.com/artisan queue:work sqs --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=forge
numprocs=8
redirect_stderr=true
#stdout_logfile=/home/forge/app.com/worker.log
stopwaitsecs=3600

# https://laravel.com/docs/9.x/queues#configuring-supervisor
